# etc/profile.d/00environ.sh
set -a
[ -f /etc/environment.base ] && source /etc/environment.base
[ -f /etc/environment ] && source /etc/environment
set +a
