# etc/profile.d/debug.sh
export MIIO_LOG_FILE=-
export MIIO_LOG_LEVEL=debug
